# repo/systemd/mail-sync.timer
[Unit]
Description=Run mail-sync every 2 Minutes

[Timer]
OnCalendar=*:00/30:00
AccuracySec=1s

[Install]
WantedBy=timers.target
